# services/startuplast.target
[Unit]
Description=Custom Target for starting up last
Requires=multi-user.target
After=multi-user.target
AllowIsolate=yes
